# deployment/backend.service
[Unit]
Description=Centinela Backend (FastAPI)
After=network.target

[Service]
Type=simple
User=centinela
WorkingDirectory=/home/centinela/centinela_vm/backend
Environment=PATH=/usr/bin:/usr/local/bin
Environment=PYTHONPATH=/home/centinela/centinela_vm/backend
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
StandardOutput=append:/home/centinela/backend.log
StandardError=append:/home/centinela/backend.log

[Install]
WantedBy=multi-user.target
